# docs/notebooks/small/windows/windows.json
{"description": "Analytics covered by the Threat Hunter Playbook Windows detection notebooks", "name": "THP Windows Analytics", "domain": "mitre-enterprise", "version": "2.2", "techniques": [{"score": 1, "techniqueID": "T1003", "metadata": [{"name": "Empire Reg Dump SAM", "value": "SD-190625133822"}, {"name": "Covenant DCSync All", "value": "SD-191027064128"}, {"name": "Covenant Mimikatz LSA Cache", "value": "SD-191205043030"}, {"name": "Covenant Mimikatz LSA Secrets", "value": "SD-191205043410"}, {"name": "Interactive Task Manager Lsass dump", "value": "SD-191027054517"}, {"name": "Remote Interactive Task Manager Lsass dump", "value": "SD-191027055035"}, {"name": "Empire Mimikatz Lsadump SAM", "value": "SD-190625103712"}, {"name": "Empire DCSync", "value": "SD-190301174830"}, {"name": "Covenant Mimikatz Logonpasswords", "value": "SD-191205033226"}, {"name": "Empire Mimikatz Export Master Key", "value": "SD-190518235535"}, {"name": "Empire Mimikatz OPTH", "value": "SD-190319131123"}, {"name": "Empire Rubeus ASKTGT", "value": "SD-190319145126"}, {"name": "Empire Mimikatz Logonpasswords", "value": "SD-190518202151"}, {"name": "Empire Rubeus ASKTGT CreateNetOnly", "value": "SD-190319151006"}, {"name": "Empire Mimikatz Extract Tickets", "value": "SD-190518230752"}, {"name": "Empire Powerdump", "value": "SD-190625132210"}, {"name": "Extended NetNTLM Downgrade", "value": "SD-191225045202"}]}, {"score": 1, "techniqueID": "T1127", "metadata": [{"name": "Covenant Grunt Msbuild", "value": "SD-191027042312"}, {"name": "Empire Invoke Msbuild", "value": "SD-190518213907"}]}, {"score": 1, "techniqueID": "T1075", "metadata": [{"name": "Empire Invoke Smbexec", "value": "SD-190518210125"}]}, {"score": 1, "techniqueID": "T1112", "metadata": [{"name": "Empire Wdigest Downgrade", "value": "SD-190518201922"}, {"name": "Empire Enabling RDP", "value": "SD-190518203650"}]}, {"score": 1, "techniqueID": "T1060", "metadata": [{"name": "Empire Userland Registry", "value": "SD-190319023812"}, {"name": "Empire Elevated Registry", "value": "SD-190518183936"}]}, {"score": 1, "techniqueID": "T1047", "metadata": [{"name": "Empire Invoke WMI Debugger", "value": "SD-190518215622"}, {"name": "Empire Invoke WMI", "value": "SD-190518214442"}, {"name": "WMIC Add User Backdoor", "value": "SD-190518231333"}]}, {"score": 1, "techniqueID": "T1069", "metadata": [{"name": "Empire Net Local Administrators Group", "value": "SD-190319020147"}, {"name": "Empire Net Domain Admins Group", "value": "SD-190518201207"}, {"name": "Empire Find Local Admin Access", "value": "SD-190518224039"}]}, {"score": 1, "techniqueID": "T1086", "metadata": [{"name": "Empire Invoke PSRemoting", "value": "SD-190518211456"}]}, {"score": 1, "techniqueID": "T1049", "metadata": [{"name": "Empire Get Local Sessions", "value": "SD-190519005224"}]}, {"score": 1, "techniqueID": "T1134", "metadata": [{"name": "Empire Invoke Runas", "value": "SD-190518204300"}]}, {"score": 1, "techniqueID": "T1084", "metadata": [{"name": "Empire Elevated WMI Subscription", "value": "SD-190518184306"}]}, {"score": 1, "techniqueID": "T1053", "metadata": [{"name": "Empire Userland Scheduled Tasks", "value": "SD-190319024742"}, {"name": "Empire Elevated Scheduled Tasks", "value": "SD-190518184109"}]}, {"score": 1, "techniqueID": "T1087", "metadata": [{"name": "Empire Net User Domain Specific", "value": "SD-190518230446"}, {"name": "Empire Net All User Domain", "value": "SD-190319021158"}, {"name": "Empire Net All Local Users", "value": "SD-190319020729"}]}, {"score": 1, "techniqueID": "T1175", "metadata": [{"name": "Empire Invoke DCOM", "value": "SD-190518211052"}]}, {"score": 1, "techniqueID": "T1055", "metadata": [{"name": "Empire DLL Injection", "value": "SD-190518221344"}, {"name": "Empire PSInject", "value": "SD-190518200432"}]}, {"score": 1, "techniqueID": "T1035", "metadata": [{"name": "Empire Invoke PsExec", "value": "SD-190518210652"}]}, {"score": 1, "techniqueID": "T1105", "metadata": [{"name": "SCM and Dll Hijacking IKEEXT", "value": "SD-190403133337"}]}, {"score": 1, "techniqueID": "T1118", "metadata": [{"name": "Covenant InstallUtil", "value": "SD-191027223020"}]}, {"score": 1, "techniqueID": "T1064", "metadata": [{"name": "Empire Launcher VBS", "value": "SD-190518182022"}]}, {"score": 1, "techniqueID": "T1222", "metadata": [{"name": "Empire DCSync ACL", "value": "SD-190301125905"}]}], "gradient": {"colors": ["#ffffff", "#66fff3"], "minValue": 0, "maxValue": 1}, "legendItems": [{"label": "Datasets researched", "color": "#66fff3"}]}